# deploy/com.novostudio.telegram-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.novostudio.telegram-mcp</string>
	<key>ProgramArguments</key>
	<array>
		<!-- Adjust node path if using a different version manager / install -->
		<string>/Users/YOUR_USER/.nvm/versions/node/v24.14.1/bin/node</string>
		<string>/Users/YOUR_USER/path/to/telegram-mcp/dist/index.js</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/YOUR_USER/path/to/telegram-mcp</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>StandardOutPath</key>
	<string>/Users/YOUR_USER/Library/Logs/telegram-mcp.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOUR_USER/Library/Logs/telegram-mcp.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/YOUR_USER</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<!-- /opt/homebrew/bin required so whisper.cpp can find ffmpeg -->
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
